# .config/autostart/ssh-add.desktop
[Desktop Entry]
Exec=ssh-add -q .ssh/id_ed25519
Name=ssh-add
Type=Application

